# systemd/user/wordle-bot.timer
[Unit]
Description=Run NYT Wordle bot daily

[Timer]
# Daily at 11:05 local time
OnCalendar=*-*-* 11:05:00
# If the timer was missed (machine asleep/off), run once when your user systemd
# instance is running again.
Persistent=true

[Install]
WantedBy=timers.target
